# addon/rootfs/etc/services.d/nginx/finish
#!/command/with-contenv bashio
# Restart the whole container if nginx dies.
bashio::log.warning "nginx exited with code $1"
exec /run/s6/basedir/bin/halt
