# addon/rootfs/etc/services.d/capturer/finish
#!/command/with-contenv bashio
bashio::log.warning "capturer exited with code $1"
exec /run/s6/basedir/bin/halt
