# addon/rootfs/etc/services.d/capturer/run
#!/command/with-contenv bashio
# Run the FastAPI capturer service.

LEVEL=$(bashio::config 'log_level')
bashio::log.info "Starting capturer FastAPI on :8765 (log level: ${LEVEL})"

cd /opt/idegis
exec python3 -m uvicorn capturer:app \
    --host 0.0.0.0 \
    --port 8765 \
    --log-level "${LEVEL}"
