# addon/rootfs/etc/services.d/capturer/run
#!/command/with-contenv bashio
# Run the aiohttp capturer service (proxy on :80, API on :8765).

LEVEL=$(bashio::config 'log_level')
bashio::log.info "Starting capturer (proxy :80, API :8765, log: ${LEVEL})"

cd /opt/idegis
exec python3 capturer.py
